# init/sqlrelay.service.in
[Unit]
Description=SQL Relay
After=network.target

[Service]
Environment=LD_LIBRARY_PATH=@libdir@
ExecStart=@bindir@/sqlr-start
Type=forking
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
